# template-java-usage/script/install/template.service
[Unit]
Description=template service
[Service]
SuccessExitStatus=143
Type=forking
ExecStart=/etc/start.sh

[Install]
WantedBy=multi-user.target